# config/docker/prod/crontab
* * * * * cd /var/www/html/ && php bin/console test